# rc.local
#!/bin/bash

/usr/sbin/nginx
/root/update.py
