# rc.local
#!/bin/bash

/usr/sbin/nginx --with-stream
/root/update.py
